# systemd/.config/systemd/user/battery-check.service
[Unit]
Description=Battery check and notification

[Service]
Type=oneshot
ExecStart=%h/battery-check.sh
