# systemd/.config/systemd/user/battery-check.timer
[Unit]
Description=Run battery check every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
Persistent=true

[Install]
WantedBy=default.target
